# .config/systemd/user/mbsync.timer
[Unit]
Description=Mbsync on a schedule

[Timer]
OnCalendar=*:*:0/30
Unit=mbsync.service

[Install]
WantedBy=default.target
