# .config/systemd/user/mbsync.service
[Unit]
Description=Checking my mail

[Service]
Type=oneshot
WorkingDirectory=%h
TimeoutStartSec=1min30s
ExecStart=/usr/bin/mbsync gmail
