# .config/systemd/user/mpd.service
[Unit]
Description=music player daemon

[Service]
Type=forking
ExecStart=/usr/bin/mpd %h/.mpd/mpd.conf
Restart=no

[Install]
WantedBy=console.target
