# .config/systemd/user/mystuff.target
[Unit]
Description=Default user target
Wants=wm.target

[Install]
Alias=default.target
