# .config/systemd/user/wm.target
[Unit]
Description=Window manager target -- xinitrc stuff
Wants=graphical.target
Requires=dbus.socket

[Install]
Alias=default.target
