# .config/systemd/user/wm@.target
[Unit]
Description=Window manager target -- xinitrc stuff
Requires=xorg@%i.target
Requires=dbus.socket
