# files/systemd/openbox@.service
[Unit]
Description=openbox window manager

[Service]
Environment=DISPLAY=:%i
ExecStart=/usr/bin/openbox

[Install]
WantedBy=wm@.target
