# roles/systemd/files/user/console.target
[Unit]
Description=Console interface target
Requires=dbus.socket
AllowIsolate=true

[Install]
Alias=default.target
