# systemd/systemd/user/xorg@.target
[Unit]
Description=Xorg target
Requires=xorg@%i.service

[Install]
WantedBy=xorg@.service
